# kafka/supervisor/image.ini
[supervisord]
nodaemon=true
user=root

[program:zookeeper]
command=/opt/zookeeper/bin/zkServer.sh start-foreground
autostart=true
autorestart=true

[program:kafka]
command=/usr/bin/start.sh
autostart=true
autorestart=true